# root/etc/services.d/fail2ban/run
#!/usr/bin/with-contenv bash
# shellcheck shell=bash

exec \
    fail2ban-client -x -vv -f -c /config/fail2ban start
